# dist/etc/systemd/system/stereopi.service
[Unit]
Description=Invoke aplay from /dev/zero at system start.

[Service]
ExecStart=sudo /home/pi/stereopi/stereopi.py

[Install]
WantedBy=multi-user.target
