# dist/etc/systemd/system/stereopi.service
[Unit]
Description=Starting the StereoPi software

[Service]
WorkingDirectory=/home/pi/stereopi
ExecStart=sudo /home/pi/stereopi/stereopi.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
